# ahtapotmys/roles/pulledpork/templates/pulledpork.timer.j2
[Unit]
Description=Pulledpork Rule Updater Scheduler Job

[Timer]
OnCalendar={{ pulledpork_update_frequency }}
Unit=pulledpork.service

[Install]
WantedBy=timers.target

